# tv-mt5-bot.service
# Systemd service file for the TradingView → MT5 bot
# Copy this to /etc/systemd/system/tv-mt5-bot.service on your AWS instance
# then run:
#   sudo systemctl daemon-reload
#   sudo systemctl enable tv-mt5-bot
#   sudo systemctl start tv-mt5-bot

[Unit]
Description=TradingView to MT5 Webhook Bot
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/TV-to-MT5-bot
EnvironmentFile=/home/ubuntu/TV-to-MT5-bot/.env
ExecStart=/usr/bin/gunicorn --workers 1 --bind 0.0.0.0:80 app:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
